# cookbooks/docker/files/etc/systemd/system/docker.service.d/environment.conf
# Load Environment
[Service]
EnvironmentFile=-/etc/default/docker
