# cookbooks/sabakan/templates/systemd.service.erb
[Unit]
Description=Sabakan Daemon
Wants=network-online.target
After=network-online.target

[Service]
PrivateTmp=true
PrivateDevices=true
ProtectSystem=full
ProtectHome=true
MountFlags=private
NoNewPrivileges=true
<%- if node[:sabakan][:config]['dhcp-bind'].match(/^.*:67$/) and node[:sabakan][:config]['http'].match(/^.*:80$/) -%>
AmbientCapabilities=CAP_NET_BIND_SERVICE CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_BIND_SERVICE CAP_NET_RAW
<%- elsif node[:sabakan][:config]['http'].match(/^.*:80$/) -%>
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
<%- end -%>
<%- case node['platform'] -%>
<%- when 'debian', 'ubuntu' -%>
EnvironmentFile=-/etc/default/sabakan
<%- when 'fedora', 'redhat', 'amazon' -%>
EnvironmentFile=-/etc/sysconfig/sabakan
<%- end -%>
User=<%= node[:sabakan][:owner] %>
Group=<%= node[:sabakan][:owner] %>
Type=simple
Restart=on-failure
RestartForceExitStatus=SIGPIPE
ExecStart=<%= "#{node[:sabakan][:prefix]}/bin/sabakan" %> -config-file <%= node[:sabakan][:file][:config][:path] %>

[Install]
WantedBy=multi-user.target
